# iot/inkbird-monitor/inkbird-monitor.service
[Unit]
Description=Inkbird IAM-T1 CO2 Sensor Monitor
Wants=network-online.target bluetooth.service
After=network-online.target bluetooth.service

[Service]
Type=simple
WorkingDirectory=/srv/selfhost/iot/inkbird-monitor
User=bleclient
Group=staff
UMask=0002
EnvironmentFile=/srv/selfhost/iot/inkbird-monitor/env
ExecStart=/srv/selfhost/iot/inkbird-monitor/bin/inkbird-monitor
Restart=always
RestartSec=10
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
